# N1/Immortalwrt/24.10/files/etc/uci-defaults/99-docker
#!/bin/sh

# 禁用/开启（disable/enable） dockerd 和 dockerman 服务
/etc/init.d/dockerd disable
/etc/init.d/dockerman disable
exit 0
